# FlockFlockDeploy/root/Library/LaunchDaemons/com.zdziarski.FlockFlockDaemon.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
        <key>UserName</key>
        <string>root</string>
        <key>GroupName</key>
        <string>wheel</string>
	<key>KeepAlive</key>
	<true/>
	<key>Nice</key>
	<integer>10</integer>
	<key>EnablePressuredExit</key>
	<false/>
	<key>ProcessType</key>
	<string>Background</string>
	<key>Label</key>
	<string>com.zdziarski.FlockFlockDaemon</string>
	<key>ProgramArguments</key>
	<array>
		<string>/Library/FlockFlock/FlockFlockDaemon</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
	<key>StandardErrorPath</key>
	<string>/Library/Logs/FlockFlockDaemon.log</string>
	<key>StandardOutPath</key>
	<string>/Library/Logs/FlockFlockDaemon.log</string>
</dict>
</plist>
